# deployment/supervisor/gopherpaw.conf
[program:gopherpaw]
; GopherPaw AI Agent Service
; Supervisor configuration file

; Command to run the application
command=/usr/local/bin/gopherpaw -c /etc/gopherpaw/config.yaml

; Directory to run the application in
directory=/var/lib/gopherpaw

; User to run as (recommended: create a dedicated user)
user=gopherpaw

; Start automatically when system boots
autostart=true

; Restart automatically if the process exits
autorestart=true

; Time to wait before restarting (in seconds)
startsecs=5

; Number of restart attempts within startretries window
startretries=3

; Time window for counting restarts (in seconds)
; If process restarts more than startretries times within this window, give up
startretries_window=60

; Redirect stderr to stdout
redirect_stderr=true

; Log file locations
stdout_logfile=/var/log/gopherpaw/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

; Environment variables
environment=
    GOPHERPAW_LLM_API_KEY="%(ENV_GOPHERPAW_LLM_API_KEY)s",
    GOPHERPAW_LLM_BASE_URL="%(ENV_GOPHERPAW_LLM_BASE_URL)s",
    GOPHERPAW_LLM_MODEL="%(ENV_GOPHERPAW_LLM_MODEL)s"

; Priority (lower values start first)
priority=999

; Stop signal (SIGTERM allows graceful shutdown)
stopsignal=TERM

; Time to wait for graceful shutdown (in seconds)
stopwaitsecs=30

; Kill the process if it doesn't stop gracefully
stopasgroup=true
killasgroup=true
